# crontab
*/15 * * * * /app/start.sh >> /app/cron.log 2>&1